# projects/Amlogic-ce/devices/Amlogic-no/packages/opengl-meson/system.d/libmali.service
[Unit]
Description=Create libMali overlay for specific AML device
After=var.mount
Before=emuelec-autostart.service

[Service]
Type=oneshot
ExecStart=/usr/sbin/libmali-overlay-setup

[Install]
WantedBy=local-fs.target
